# root/etc/s6-overlay/s6-rc.d/svc-msmtpd/run
#!/usr/bin/execlineb -P
with-contenv
#echo "[svc-msmtpd] Redirecting file descriptors..."
fdmove -c 2 1
# fdmove 1 3
s6-setuidgid abc 
/app/start-msmtpd