# root/etc/s6-overlay/s6-rc.d/init-timezone/run
#!/usr/bin/with-contenv bash
# spellcheck shell=bash

TZ=${TZ:-UTC}

echo "Setting timezone to ${TZ}..."
ln -snf /usr/share/zoneinfo/${TZ} /etc/localtime
echo ${TZ} > /etc/timezone