# root/etc/s6-overlay/s6-rc.d/init-folders/run
#!/usr/bin/with-contenv bash
# spellcheck shell=bash

# make folders
echo "[init-folders] Ensure required folder structure exists..."
mkdir -p \
  /config/log/msmtpd