# root/etc/s6-overlay/s6-rc.d/init-permissions/run
#!/usr/bin/with-contenv bash
# spellcheck shell=sh

echo "[init-permissions] Setting permissions..."
lsiown -R abc:abc \
  /config/msmtp.conf \
  /config/log

chmod 600 /config/msmtp.conf